# supervisord.conf
[supervisord]
nodaemon=true

[program:http-server]
command=http-server

[program:nodejs]
command=nodejs server/app.js